# deploy/managrams.timer.tmpl
[Unit]
Description=Process managrams every minute

[Timer]
OnActiveSec=1
OnUnitActiveSec=1m
Persistent=true

[Install]
WantedBy=timers.target