# deploy/sync.timer.tmpl
[Unit]
Description=Do a full sync every 4 hours

[Timer]
OnActiveSec=15
OnUnitActiveSec=4h
Persistent=true

[Install]
WantedBy=timers.target